# imageapp.service
[Unit]
Description=AI Image Studio
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/exedev/imageapp
ExecStart=/home/exedev/imageapp/.venv/bin/python main.py
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
