# retouch.service
[Unit]
Description=Retouch - AI Image Studio
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/exedev/retouch
ExecStart=/home/exedev/retouch/.venv/bin/python main.py
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
